# deploy/systemd/aethervault-self-improve.timer
[Unit]
Description=AetherVault Self-Improvement Timer (every 6h)

[Timer]
OnCalendar=*-*-* 00/6:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
