# deploy/systemd/aethervault.service
[Unit]
Description=AetherVault Telegram Bridge (AetherVault runtime)
After=network.target ollama.service
Wants=ollama.service
OnFailure=aethervault-rollback.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/.aethervault
EnvironmentFile=/root/.aethervault/.env
# Apply patches before starting (survives npm updates)
ExecStartPre=/root/.aethervault/patches/apply-patches.sh
ExecStart=/usr/local/bin/aethervault bridge telegram --log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
